# deploy/worker.service
[Unit]
Description=CarbonTracker Celery Worker
After=network-online.target
Wants=network-online.target

[Service]
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/carbon-tracker
EnvironmentFile=/etc/default/carbon-app
Environment=PYTHONPATH=/opt/carbon-tracker/app.py
ExecStart=/opt/carbon-tracker/.venv/bin/celery \
  -A celery_app:celery worker \
  --queues=carbon-jobs \
  --include=tasks \
  --loglevel=info
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target